# systemd/varnish-reloader.path
# block-excessive-ua.path
[Unit]
Description=Monitor the blocked file for changes

[Path]
PathChanged=/etc/varnish/blocked-ua.vcl
Unit=varnish-reloader.service

[Install]
WantedBy=multi-user.target
